# ops/local-staging/systemd/veud-staging-mal-trending.timer
[Unit]
Description=Refresh Veud staging MAL momentum every six hours

[Timer]
OnCalendar=*-*-* 00,06,12,18:20:00
RandomizedDelaySec=10m
Persistent=true
Unit=veud-staging-mal-trending.service

[Install]
WantedBy=timers.target
