# ops/local-staging/systemd/veud-staging-mal-inventory.timer
[Unit]
Description=Refresh Veud staging MAL ranking inventory daily

[Timer]
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=15m
Persistent=true
Unit=veud-staging-mal-inventory.service

[Install]
WantedBy=timers.target
