# ops/local-staging/systemd/veud-staging-app.service
[Unit]
Description=Veud isolated staging application
After=veud-staging-postgres.service
Requires=veud-staging-postgres.service
ConditionPathIsMountPoint=/media/sde

[Service]
Type=simple
ExecStart=%h/Programs/Veud/ops/local-staging/run-app.sh
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
