# ops/local-staging/systemd/veud-staging-catalog-backup.timer
[Unit]
Description=Daily Veud staging catalog backup and restore verification

[Timer]
OnCalendar=*-*-* 04:30:00
RandomizedDelaySec=15m
Persistent=true
Unit=veud-staging-catalog-backup.service

[Install]
WantedBy=timers.target
