# ops/local-production/systemd/veud-production-tmdb-hydration.timer
[Unit]
Description=Resume Veud production TMDB catalog hydration daily

[Timer]
OnBootSec=3min
OnUnitInactiveSec=1d
Persistent=true
Unit=veud-production-tmdb-hydration.service

[Install]
WantedBy=timers.target
